# templates/velveta-cleanup.timer.j2
[Unit]
Description=Run Velveta cleanup on schedule

[Timer]
OnCalendar={{ velveta_cleanup.timer_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
